# ansible/roles/oasys-ords/templates/tomcat-daily-restart.timer.j2
[Unit]
Description=Timer for tomcat nightly service restart

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
