# ansible/roles/oasys-ords/templates/tomcat-daily-restart.service.j2
[Unit]
Description=Restart tomcat nightly

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart tomcat.service
