# ansible/roles/collectd-service-metrics/files/oracle-batch-status.service
[Service]
ExecStart=/opt/collectd/oracle_batch_status.sh
