# ansible/roles/collectd-service-metrics/files/oracle-db-connected.service
[Service]
ExecStart=/opt/collectd/oracle_db_connected.sh
