# ansible/roles/autoscale-group-hooks/templates/autoscaling-lifecycle-hook.cron.j2
@reboot root /usr/local/bin/autoscaling-lifecycle-{{ lifecycle_hook_name }}.sh
